# units/healthmon.service
# /etc/systemd/system/healthmon.service
[Unit]
Description=Server Health Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/healthmon.sh --silent --email --slack

[Install]
WantedBy=multi-user.target 
